# myproject.conf
description "dharmana upstart"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid root
setgid www-data

env PATH=/home/ubuntu/myproject/myproject/bin
chdir /home/ubuntu/myproject
exec gunicorn --workers 3 --bind unix:myproject.sock -m 007 wsgi
